# pishutdownbutton_1.0-2/DEBIAN/postinst
# postinst
sudo pip3 install --upgrade RPi.GPIO
systemctl enable pishutdownbutton
systemctl start pishutdownbutton
